# kata-agent.service.in
#
# Copyright (c) 2018-2019 Intel Corporation
#
# SPDX-License-Identifier: Apache-2.0
#

[Unit]
Description=Kata Containers Agent
Documentation=https://github.com/kata-containers/agent
Wants=kata-containers.target

[Service]
# Send agent output to tty to allow capture debug logs
# from a VM serial port
StandardOutput=tty
Type=simple
ExecStart=@bindir@/@kata-agent@
LimitNOFILE=infinity
ExecStop=/bin/sync
